# panel/mihomo-gateway-api.service
[Unit]
Description=Mihomo Gateway management API
After=network-online.target mihomo.service
Wants=network-online.target

[Service]
Type=simple
Environment=PANEL_API_HOST=127.0.0.1
Environment=MIHOMO_CONFIG=/etc/mihomo/config.yaml
Environment=PROTECTED_PROVIDERS=
WorkingDirectory=__INSTALL_ROOT__/panel
ExecStart=/bin/bash -c 'set -a; source /root/mihomo-gateway/env; set +a; exec /usr/bin/python3 __INSTALL_ROOT__/panel/app.py'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
